# config/systemd/restore-vesktop-css.service
[Unit]
Description=Restore Vesktop quickCss.css after updates
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/%u/.local/bin/restore-vesktop-css.sh
RemainAfterExit=yes

[Install]
WantedBy=default.target
