# config/systemd/restore-vesktop-css.path
[Unit]
Description=Monitor Vesktop configuration for updates to restore quickCss.css

[Path]
PathChanged=$VESKTOP_CSS_FILE
Unit=restore-vesktop-css.service

[Install]
WantedBy=default.target
