# config/systemd/restore-vesktop-css.path
[Unit]
Description=Monitor Vesktop configuration for updates to restore quickCss.css

[Path]
PathChanged=/home/%u/.var/app/dev.vencord.Vesktop/config/vesktop/settings/quickCss.css
Unit=restore-vesktop-css.service

[Install]
WantedBy=default.target
